# streamer/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:dbus]
command=/usr/bin/dbus-daemon --system --nofork --nopidfile
autostart=true
autorestart=true
priority=1

[program:avahi]
command=/usr/sbin/avahi-daemon --no-drop-root
autostart=true
autorestart=true
priority=5

[program:nqptp]
command=/usr/local/bin/nqptp
autostart=true
autorestart=true
priority=10

[program:pulseaudio]
command=/usr/bin/pulseaudio --daemonize=no --disallow-exit --exit-idle-time=-1
environment=DBUS_SYSTEM_BUS_ADDRESS="unix:path=/var/run/dbus/system_bus_socket",DBUS_SESSION_BUS_ADDRESS="unix:path=/var/run/dbus/system_bus_socket"
autostart=true
autorestart=true
priority=10
startretries=3
startsecs=2
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:shairport-sync]
command=/bin/bash -c "while [ ! -S /var/run/pulse/pulseaudio.socket ]; do sleep 0.5; echo 'Waiting for pulseaudio...'; done; exec /usr/local/bin/shairport-sync -a \"%(ENV_SOUND_DEVICE_NAME)s\" --output pulseaudio -- -s unix:/var/run/pulse/pulseaudio.socket"
autostart=true
autorestart=true
priority=20
startretries=3
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
